# deploy/conary-test.service
[Unit]
Description=Conary Test Infrastructure Server
After=network.target

[Service]
Type=simple
User=peter
WorkingDirectory=/home/peter/Conary
EnvironmentFile=-/etc/conary-test/env
Environment=DOCKER_HOST=unix:///run/user/1000/podman/podman.sock
ExecStart=/home/peter/Conary/target/debug/conary-test serve --port 9090 --token ${CONARY_TEST_TOKEN}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
