# apps/conary/tests/fixtures/phase4-daily-driver-corpus/stage/usr/lib/systemd/system/phase4-corpus.service
[Unit]
Description=Phase 4 corpus service

[Service]
Type=oneshot
ExecStart=/usr/bin/phase4-corpus

[Install]
WantedBy=multi-user.target
